# .config/systemd/user/kdeconnect-mic-mute.service
[Unit]
Description=Mute microphone on KDE Connect phone calls
After=graphical-session.target kdeconnect.target
PartOf=graphical-session.target

[Service]
ExecStart=%h/bin/kdeconnect-mic-mute
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=graphical-session.target
